# recipes-rpi-gui/pyscript/files/pyscript.service
[Unit]
Description=Run python script at startup

[Service]
ExecStart=/usr/bin/python3 /usr/share/script.py
Restart=always
User=root
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
